# templates/runit/log/run
#!/bin/sh

exec chpst -ulogger svlogd -tt ./main
